# noaaSourcePaths.cron
06 11 * * * /Users/Neville/anaconda/bin/python /Users/Neville/software/PolarExplorerWebApp/src/getLatestNOAAImages.py /Users/Neville/software/PolarExplorerWebApp
